# html2x.conf
[supervisord]
nodaemon=true

[program:html2x]
directory=/home/app;
command=./html2x;
autostart=true;
autorestart=true;
user=root;
redirect_stderr=true;
stdout_logfile=/var/log/supervisor/html2x.log;
